# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:log.io-server]
command=/usr/bin/log.io-server
autorestart=true

[program:log.io-harvester]
command=/usr/bin/log.io-harvester
autorestart=true
